# autossh/rootfs/etc/services.d/autossh/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Autossh
# Runs Autossh script
# ==============================================================================


bashio::log.info "Starting Autossh service..."

chmod +x /usr/bin/autossh.sh
exec /usr/bin/autossh.sh
